# launchd/ai.openclaw.gmail_sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gmail_sentinel</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>${HOME}/.openclaw/bin/gmail_heartbeat_sentinel.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/openclaw/gmail_sentinel.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/openclaw/gmail_sentinel.err.log</string>
  </dict>
</plist>
